# arego-vite.service
[Unit]
Description=Arego Chat — Vite Dev Server (HTTP :5173, hinter Nginx)
Documentation=https://github.com/Aregoland/arego-chat
After=network.target arego-signaling.service
Wants=arego-signaling.service

[Service]
Type=simple
Restart=always
RestartSec=5
WorkingDirectory=/root/Aregoland

# Vite lauscht auf 127.0.0.1:5173, Nginx macht SSL-Terminierung
Environment=NODE_ENV=development
ExecStart=/usr/bin/pnpm dev

# Sanfter Stop: SIGINT damit Vite sauber beendet
KillSignal=SIGINT
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
