# arego-vite.service
[Unit]
Description=Arego Chat — Vite Dev Server (HTTPS :443)
Documentation=https://github.com/Aregoland/arego-chat
After=network.target arego-signaling.service
Wants=arego-signaling.service

[Service]
Type=simple
Restart=always
RestartSec=5
WorkingDirectory=/root/Aregoland

# PORT 443 braucht root oder CAP_NET_BIND_SERVICE
# Läuft als root (Standard wenn kein User= gesetzt)
Environment=NODE_ENV=development
ExecStart=/usr/bin/pnpm dev

# Sanfter Stop: SIGINT damit Vite sauber beendet
KillSignal=SIGINT
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
